# Firmware/OpenWrt/packages/etherdrive-webui/files/etc/init.d/etherdrive
#!/bin/sh /etc/rc.common

START=99
STOP=10
USE_PROCD=0

boot() {
	start
}

start() {
	[ -x /usr/bin/etherdrive-storage ] && /usr/bin/etherdrive-storage ensure-mounted >/dev/null 2>&1 || true
	[ -x /usr/bin/etherdrive-mode ] && /usr/bin/etherdrive-mode apply >/dev/null 2>&1 || true
}

stop() {
	[ -x /etc/init.d/samba ] && /etc/init.d/samba stop >/dev/null 2>&1 || true
	[ -x /etc/init.d/udpfsd ] && /etc/init.d/udpfsd stop >/dev/null 2>&1 || true
}

reload() {
	start
}
